# resources/scripts-deb/prerm
#!/bin/bash
/bin/systemctl stop stalwart-imap
/bin/systemctl disable stalwart-imap
deluser stalwart-imap
